# packaging/systemd/zen-history.timer
[Unit]
Description=Record Zen pressure history snapshots periodically
Documentation=https://github.com/dodge1218/zen

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
